# programs/ziti-edge-tunnel/package/systemd/ziti-edge-tunnel.service.in
[Unit]
Description=Ziti Edge Tunnel
After=network-online.target

[Service]
Type=simple
# default values
Environment="ZITI_IDENTITY_DIR=/opt/openziti/etc/identities" "ZITI_DNS_IP_RANGE=100.64.0.1/10" "ZITI_VERBOSE=2"
# optional override values
EnvironmentFile=-@CPACK_ETC_DIR@/@SYSTEMD_SERVICE_NAME@.env
User=ziti
UMask=0007
AmbientCapabilities=CAP_NET_ADMIN
ExecStartPre=@CPACK_BIN_DIR@/@SYSTEMD_SERVICE_NAME@.sh
ExecStart=@CPACK_BIN_DIR@/@SYSTEMD_SERVICE_NAME@ run --verbose=${ZITI_VERBOSE} --dns-ip-range=${ZITI_DNS_IP_RANGE} --identity-dir=${ZITI_IDENTITY_DIR}
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
