# programs/ziti-edge-tunnel/package/deb/postrm.in
# Copied from https://github.com/Debian/debhelper/blob/master/autoscripts/postrm-systemd
if [ "$1" = "remove" ]; then
    if [ -x "/usr/bin/deb-systemd-helper" ]; then
        deb-systemd-helper mask @SYSTEMD_UNIT_FILE_NAME@ >/dev/null || true
    fi
fi

if [ "$1" = "purge" ]; then
    if [ -x "/usr/bin/deb-systemd-helper" ]; then
        deb-systemd-helper purge @SYSTEMD_UNIT_FILE_NAME@ >/dev/null || true
        deb-systemd-helper unmask @SYSTEMD_UNIT_FILE_NAME@ >/dev/null || true
    fi
fi
# End copied section

# delete the symlink to the executable that was created by the postinst scriptlet
if [ -L /usr/bin/@CPACK_PACKAGE_NAME@ ]; then
    unlink /usr/bin/@CPACK_PACKAGE_NAME@
fi

# delete the regular file of the service unit that was copied by the postinst scriptlet
if [ -e @SYSTEMD_UNIT_DIR@/@SYSTEMD_UNIT_FILE_NAME@ ]; then
    unlink @SYSTEMD_UNIT_DIR@/@SYSTEMD_UNIT_FILE_NAME@
fi

# Copied from https://github.com/Debian/debhelper/blob/master/autoscripts/postrm-systemd-reload-only
if [ -d /run/systemd/system ]; then
    systemctl --system daemon-reload >/dev/null || true
fi
# End copied section

if [ -e /usr/share/debconf/confmodule ]; then
    # Source debconf library.
    . /usr/share/debconf/confmodule

    db_get ziti_edge_tunnel/install_pkla
    if [ "$RET" = "true" ]; then
      rm -f "/var/lib/polkit-1/localauthority/10-vendor.d/@ZITI_POLKIT_PKLA_FILE@"
      db_set ziti_edge_tunnel/install_pkla false
    fi

    if [ "$1" = "purge" ]; then
      # Remove changes to the db.
      db_purge
    fi
fi

###
### Unlikely all files and directories were removed.
### Some of the unremoved files and directories are likely owned by `ziti' user or group.
### Consequently, don't remove 'ziti` user and group as it will strand the files
###
#if [ "$1" = "purge" ]; then
#  if command -v deluser >/dev/null; then
#    deluser --quiet ziti >/dev/null || true
#    delgroup --group --quiet ziti >/dev/null || true
#  fi
#fi
