# programs/ziti-edge-tunnel/package/deb/postinst.in
ln -sf @CPACK_BIN_DIR@/@CPACK_PACKAGE_NAME@ /usr/bin/@CPACK_PACKAGE_NAME@
ln -sf @CPACK_SHARE_DIR@/@SYSTEMD_UNIT_FILE_NAME@ /etc/systemd/system/@SYSTEMD_UNIT_FILE_NAME@

# Add user `ziti'
if [ "$1" = "configure" ]; then
  # create user / group
  echo 'u ziti - "openziti user" "/var/lib/ziti"' | \
    systemd-sysusers --replace=/usr/lib/sysusers.d/ziti-edge-tunnel.conf -

  # update permissions
  chown -R ziti:ziti "@ZITI_STATE_DIR@"
  chmod -R u=rwX,g=rwX,o= "@ZITI_STATE_DIR@"

  chown root:ziti "@ZITI_IDENTITY_DIR@"
  chmod 0770 "@ZITI_IDENTITY_DIR@"
  find "@ZITI_IDENTITY_DIR@" -maxdepth 1 -name "*.json" -type f -exec chown ziti:ziti "{}" + -exec chmod 0400 "{}" +

  # install PolicyKit policies based on policy directory existance
  ucf "@CPACK_SHARE_DIR@/@ZITI_POLKIT_RULES_FILE@.sample" "/usr/share/polkit-1/rules.d/@ZITI_POLKIT_RULES_FILE@" 2>/dev/null || true
  ucf "@CPACK_SHARE_DIR@/@ZITI_POLKIT_PKLA_FILE@.sample" "/var/lib/polkit-1/localauthority/10-vendor.d/@ZITI_POLKIT_PKLA_FILE@" 2>/dev/null || true
fi

# Copied from https://github.com/Debian/debhelper/blob/master/autoscripts/postinst-systemd-dont-enable
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if deb-systemd-helper debian-installed @SYSTEMD_UNIT_FILE_NAME@; then
		# This will only remove masks created by d-s-h on package removal.
		deb-systemd-helper unmask @SYSTEMD_UNIT_FILE_NAME@ >/dev/null || true

		if deb-systemd-helper --quiet was-enabled @SYSTEMD_UNIT_FILE_NAME@; then
			# Create new symlinks, if any.
			deb-systemd-helper enable @SYSTEMD_UNIT_FILE_NAME@ >/dev/null || true
		fi
	fi

	# Update the statefile to add new symlinks (if any), which need to be cleaned
	# up on purge. Also remove old symlinks.
	deb-systemd-helper update-state @SYSTEMD_UNIT_FILE_NAME@ >/dev/null || true
fi
# End copied section

if [ "$1" = "configure" ]; then
    ssize=$(tput cols)
    printf '\n'
    printf %"$ssize"s | tr " " "-"
    echo "@SYSTEMD_SERVICE_NAME@ was installed..."
    echo "First install an OpenZiti identity or enroll token in: @ZITI_IDENTITY_DIR@"
    echo "then start or restart this systemd service unit."
    printf %"$ssize"s | tr " " "-"
    printf '\n'
fi
